# Lab6-SoftwareRouter/supervisord_conf.d/bgpd.conf
[program:bgpd]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=bgpd -A 127.0.0.1 -f /etc/quagga/bgpd.conf

