# Lab6-SoftwareRouter/supervisord_conf.d/zebra.conf
[program:zebra]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=zebra -A 127.0.0.1 -f /etc/quagga/zebra.conf

